# rootfs/etc/s6-overlay/s6-rc.d/fava-log-prepare/up
if { mkdir -p /var/log/fava }
if { chown nobody:nogroup /var/log/fava }
chmod 02755 /var/log/fava
