# rootfs/etc/s6-overlay/s6-rc.d/fava/run
#!/command/execlineb -P

multisubstitute
{
  import -u BEANCOUNT_FILE
}

s6-setuidgid beancount
fava ${BEANCOUNT_FILE}
