# rootfs/etc/s6-overlay/s6-rc.d/fava/run
#!/usr/bin/execlineb -P

with-contenv

s6-setuidgid fava
elglob -s -0 beanfiles /data/*.bean
elglob -s -0 beancountfiles /data/*.beancount
fava ${beanfiles} ${beancountfiles}
